# packages/front_debian_src/postrm
#!/bin/bash
set -e

update-rc.d -f front_d remove